# systemd/pacwall.service
[Unit]
Description=Pacwall
ConditionUser=!root

[Service]
Type=oneshot
ExecStart=/usr/bin/pacwall
