# systemd/pacwall.path
[Unit]
Description=Run pacwall on package update/removal/installation

[Path]
PathChanged=/var/lib/pacwall/systemd/trigger

[Install]
WantedBy=default.target
