# systemd/pacwall.service
[Unit]
Description=Pacwall wallpaper generator
ConditionUser=!root

[Service]
Type=oneshot
ExecStart=/usr/bin/pacwall
