# systemd/pacwall.service
[Unit]
Description=Pacwall wallpaper generator
ConditionUser=!root

[Service]
ExecStart=/usr/bin/pacwall
Type=oneshot
Restart=no
