# systemd/pacwall-watch-packages.path
[Unit]
Description=Run pacwall on package upgrade/removal/installation

[Path]
Unit=pacwall@u.service
PathChanged=/var/lib/pacwall/systemd/trigger

[Install]
WantedBy=default.target
